# root/etc/services.d/ozwcp/run
#!/usr/bin/with-contenv bash
cd /config/ozwcp/ || exit
exec \
        s6-setuidgid abc \
		/config/ozwcp/ozwcp -p 8888
